# roles/vmalert/templates/vmalert.service
[Unit]
Description=vmalert - VictoriaMetrics Alerting and Recording Rules
Documentation=https://docs.victoriametrics.com/vmalert.html
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ vmalert_user }}
Group={{ vmalert_group }}
ExecStart=/usr/bin/vmalert \
    -httpListenAddr=:{{ vmalert_port }} \
    -datasource.url={{ victoria_metrics_url }} \
    -remoteWrite.url={{ victoria_metrics_write_url }} \
    -remoteRead.url={{ victoria_metrics_url }}{% if vmalert_alertmanager_enabled|bool %} \
    -notifier.url={{ alertmanager_url }}{% endif %} \
    -rule={{ vmalert_rules_dir }}/*.yml \
    -evaluationInterval={{ vmalert_evaluation_interval }} \
    -loggerLevel={{ vmalert_log_level }} \
    -rule.updateEntriesLimit={{ vmalert_rule_update_entries_limit }}{% if vmalert_external_labels %}{% for key, value in vmalert_external_labels.items() %} \
    -external.label={{ key }}={{ value }}{% endfor %}{% endif %}

ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
LimitNPROC=32768

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vmalert_data_dir }}
ReadOnlyPaths={{ vmalert_config_dir }} /etc/prometheus
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target