# roles/vlagent/templates/vlagent.service
[Unit]
Description=vlagent - VictoriaLogs Log Collection Agent
Documentation=https://docs.victoriametrics.com/VictoriaLogs/data_ingestion.html#vlagent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=vlagent
Group=vlagent
WorkingDirectory={{ vlagent_data_dir }}
ExecStart=/usr/bin/vlagent \
    -httpListenAddr=:{{ vlagent_port }} \
    -loggerLevel={{ vlagent_log_level }} \
    -remoteWrite.tmpDataPath={{ vlagent_data_dir }}/remotewrite-data \
{% for host in groups['infra'] %}
    -remoteWrite.url=http://{{ host }}:3100/internal/insert \
{% endfor %}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vlagent_data_dir }}

PrivateTmp=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target