# roles/vlagent/templates/vlagent.service
[Unit]
Description=vlagent - VictoriaLogs Log Collection Agent
Documentation=https://docs.victoriametrics.com/VictoriaLogs/data_ingestion.html#vlagent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=vlagent
Group=vlagent
ExecStart=/usr/bin/vlagent \
    -config.file={{ vlagent_config_dir }}/vlagent.yml \
    -httpListenAddr=:{{ vlagent_port }} \
    -loggerLevel={{ vlagent_log_level }}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vlagent_data_dir }} /var/log{% if docker_enabled|default(false) %} /var/lib/docker/containers{% endif %}

PrivateTmp=true
{% if docker_enabled|default(false) and vlagent_docker_group_exists|default(false) %}
SupplementaryGroups=docker
{% endif %}
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target